# docker-build_moinmoin/moin.conf
description "moin uwsgi service"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /srv/www/moin/wiki/
exec uwsgi -x /srv/www/moin/wiki/uwsgi.xml
respawn
